# scripts/spored.service
[Unit]
Description=Spored Agent - Instance self-monitoring and auto-termination
Documentation=https://github.com/spore-host/spawn
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/spored
# on-failure (not always) prevents restart attempts during graceful shutdown
Restart=on-failure
RestartSec=10
# Give spored time to deregister DNS and clean up before SIGKILL
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
# NOTE: do NOT set PrivateTmp=true. spored watches the host-side completion
# file (default /tmp/SPAWN_COMPLETE) written by users, `spawn connect`, and
# nf-spawn. PrivateTmp gives the daemon an isolated /tmp where it can never
# see that file, silently disabling --on-complete and --pre-stop (#66).

# Resource limits
LimitNOFILE=8192

[Install]
WantedBy=multi-user.target
